# geordi/supervisor.conf
[supervisord]
nodaemon=true

[program:geordi]
user=root
command=/geordi/node_modules/.bin/slc run .
directory=/geordi
redirect_stderr=true
autorestart=true
stdout_logfile=/geordi/log/stdout.log
