# geordi/supervisor.conf
[supervisord]
nodaemon=true

[program:geordi]
user=root
command=/app/node_modules/.bin/slc run .
directory=/app
redirect_stderr=true
autorestart=true
stdout_logfile=/app/log/stdout.log
