# geordi/supervisor.conf
[supervisord]
nodaemon=true

[program:geordi]
user=root
command=slc run .
directory=/app
redirect_stderr=true
autorestart=true
stdout_logfile=/app/log/stdout.log
