# API/supervisor.conf
[supervisord]
nodaemon=true

[program:api]
user=root
command=node geordi.js 80
directory=/app
redirect_stderr=true
autorestart=true
stdout_logfile=/app/log/stdout.log
